# packaging/systemd/openclaw-portable-bridge-broker.service.example
[Unit]
Description=OpenClaw Portable Bridge pairing broker
After=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/openclaw-portable-bridge/broker.env
ExecStart=%h/.local/lib/openclaw-portable-bridge/pairing-broker -listen 127.0.0.1:17443 -audit %h/.local/state/openclaw-portable-bridge/audit.jsonl
Restart=on-failure
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/openclaw-portable-bridge

[Install]
WantedBy=default.target

